# backend/docker/supervisord/syslog.conf
[group:syslog]
programs=syslogd
priority=10

[program:syslogd]
command = /opt/docker/bin/service.d/syslog-ng.sh
process_name=%(program_name)s
autostart = false
autorestart = true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
startretries=20
